# cron/fetchprice.cron
*/15 * * * * /usr/bin/php /var/www/sumblockhalf.com/getprice.php